# com.tracker.summary.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tracker.summary</string>

    <key>ServiceDescription</key>
    <string>Web Activity Summarizer</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON_PATH__</string>
        <string>__PROJECT_DIR__/summarizer.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <!-- Start on login and keep alive — APScheduler handles the schedule internally -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Schedule + API keys (filled in by install.sh) -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>OPENAI_API_KEY</key>
        <string>__OPENAI_API_KEY__</string>
        <key>MAILERSEND_API_KEY</key>
        <string>__MAILERSEND_API_KEY__</string>
        <key>MAILERSEND_FROM</key>
        <string>__MAILERSEND_FROM__</string>
        <key>MAILERSEND_TO</key>
        <string>__MAILERSEND_TO__</string>
        <key>OPENAI_MODEL</key>
        <string>__OPENAI_MODEL__</string>
        <key>SUMMARY_SCHEDULE</key>
        <string>__SUMMARY_SCHEDULE__</string>
        <key>SUMMARY_SCHEDULE_HOUR</key>
        <string>__SUMMARY_SCHEDULE_HOUR__</string>
        <key>SUMMARY_SCHEDULE_MINUTE</key>
        <string>__SUMMARY_SCHEDULE_MINUTE__</string>
        <key>SUMMARY_SCHEDULE_WEEKDAY</key>
        <string>__SUMMARY_SCHEDULE_WEEKDAY__</string>
        <key>SUMMARY_SCHEDULE_DAY</key>
        <string>__SUMMARY_SCHEDULE_DAY__</string>
        <key>SUMMARY_SCHEDULE_INTERVAL_MINUTES</key>
        <string>__SUMMARY_SCHEDULE_INTERVAL_MINUTES__</string>
    </dict>

    <!-- Python writes rotating logs internally; only capture stderr from launchd -->
    <key>StandardOutPath</key>
    <string>/dev/null</string>
    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/summarizer_stderr.log</string>
</dict>
</plist>
